# static/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

umask 0002

echo "Starting Subsonic"
exec s6-setuidgid bytesized /app/subsonic.sh
echo "Sonarr shut down"
